# 3n1-0.1_amd64/DEBIAN/postinst
#!/usr/bin/bash

systemctl daemon-reload

systemctl enable 3n1
systemctl start 3n1

systemctl enable 3n1shutdown
systemctl start 3n1shutdown